# work/peer1/kickstart/install/var/www/kickstart/bin/mac_cleaner.cron
# m h  dom mon dow   command
0 3 * * sun /exports/kickstart/bin/mac_cleaner.pl >> /var/log/mac_cleaner.log 2>&1

